# scripts/com.atom.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.atom.agent</string>
	<key>ProgramArguments</key>
	<array>
		<string>@@@PYTHON3@@@</string>
		<string>@@@ATOM_REPO@@@/main.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>@@@ATOM_REPO@@@</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>15</integer>
	<key>StandardOutPath</key>
	<string>@@@ATOM_REPO@@@/logs/launchagent.stdout.log</string>
	<key>StandardErrorPath</key>
	<string>@@@ATOM_REPO@@@/logs/launchagent.stderr.log</string>
</dict>
</plist>
